# supervisord.conf
[supervisord]
nodaemon=true

[program:v2ray]
command=/usr/bin/v2ray -c /etc/v2ray/config.json

[program:v2ray-mu]
command=/usr/bin/v2ray-mu